# scripts/heartbeat-server.service
[Unit]
Description=Runs Heatbeat Beat
Wants=heartbeat-server.timer

[Service]
Type=oneshot
ExecStart=/home/heartbeat/heartbeat/start.sh

[Install]
WantedBy=multi-user.target