# scripts/heartbeat-client.timer
[Unit]
Description=Runs Heatbeat Beat
Requires=heartbeat-client.service

[Timer]
Unit=heartbeat-client.service
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target